# per-user/mediagoblin/supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=sudo -u postgres  /usr/lib/postgresql/9.1/bin/postgres -D /var/lib/postgresql/9.1/main -c config_file=/etc/postgresql/9.1/main/postgresql.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=false
autostart=true
redirect_stderr=true

[program:mediagoblin]
command=/srv/mediagoblin/run_app.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=false
autostart=true
redirect_stderr=true
