# per-user/lamp-git/backup.cron
42 * * * * root cd /data; mysqldump --all-databases > dump.sql; git add *; git commit -am"backup `date`"

